# arkime6/roles/arkime_recorder/templates/capture-nic.service.j2
{{ ansible_managed | comment }}
# Oneshot unit that re-applies Arkime capture NIC tuning on every boot:
# promiscuous mode on plus hardware offloads (gro/lro/tso/gso) off, so capture
# always sees real wire frames regardless of NIC driver defaults.
[Unit]
Description=FPC Arkime capture NIC tuning (promisc + offloads off)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
{% for iface in capture_interfaces %}
{% if capture_promisc | bool %}
ExecStart=/usr/sbin/ip link set {{ iface }} promisc on
{% endif %}
{% if capture_disable_offloads | bool %}
ExecStart=-/usr/sbin/ethtool -K {{ iface }} gro off lro off tso off gso off
{% endif %}
{% endfor %}

[Install]
WantedBy=multi-user.target
